# Linux/Dotfiles/usr/lib/systemd/system/lowpower.service
[Unit]
Description=Low power notify and suspend

[Service]
ExecStart=/usr/bin/lowpower.sh
