# Arch/Dotfiles/etc/udev/rules.d/99-hdmi_sound.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/hdmi_sound_toggle.sh"
