# Gentoo/etc/systemd/system/hdmi_sound_toggle.service
[Unit]
Description=hdmi sound hotplug

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/bin/hdmi_sound_toggle.sh

[Install]
WantedBy=multi-user.target
