# Linux/Dotfiles/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin electrux --noclear %I $TERM
